# gar-minimyth/script/meta/minimyth/dirs/etc/inittab
::sysinit:/etc/rc

vc/1::respawn:-/bin/sh 
vc/2::respawn:-/bin/sh 
vc/3::respawn:-/bin/sh 
vc/4::respawn:-/bin/sh 

::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -af 

