# trunk/gar-minimyth/script/meta/minimyth/dirs/etc/inittab
::sysinit:/etc/rc

tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh
tty4::askfirst:-/bin/sh
tty5::askfirst:-/bin/sh

::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -af 
