# trunk/gar-minimyth/script/meta/minimyth/files/source/rootfs/lib/udev/rules.d/10-minimyth-asound.rules
ACTION!="add|change", GOTO="end"

SUBSYSTEM=="sound", RUN+="/lib/udev/mm_asound"

LABEL="end"
