# trunk/gar-minimyth/script/meta/minimyth/dirs/udev/rules.d/03-minimyth-local.rules
# Handle the MiniMyth local boot/conf device.
#   Identify it based on it being a block disk drive/partition with the label 'minimyth'.
#   Change its device name to 'disk-minimyth' so that the MiniMyth patched MythTV will ignore it.
#   Mount is as '/minimyth' so that MiniMyth can find it.


SUBSYSTEM=="block", KERNEL=="hd[a-z]*|sd[a-z]*|sr[0-9]*", \
	IMPORT{program}="/lib/udev/vol_id --export $tempnode"

SUBSYSTEM=="block", KERNEL=="hd[a-z]*|sd[a-z]*|sr[0-9]*", ENV{ID_FS_LABEL_SAFE}=="minimyth", \
	NAME="disk-minimyth", SYMLINK+="%k"

NAME=="disk-minimyth", ACTION=="add", \
	RUN+="/lib/udev/mm_mount add $env{DEVNAME} 1 /minimyth $env{ID_FS_TYPE}"
NAME=="disk-minimyth", ACTION=="delete", \
	RUN+="/lib/udev/mm_mount remove $env{DEVNAME}"
