# trunk/gar-minimyth/script/meta/minimyth/dirs/udev/rules.d/04-minimyth-mythtv.rules
SUBSYSTEM=="block", RUN+="/lib/udev/mm_mythtv"
