# gar-minimyth/script/meta/minimyth/files/source/rootfs/usr/lib/udev/rules.d/10-minimyth-asound.rules
ACTION!="add|change", GOTO="end"

SUBSYSTEM=="sound", RUN+="mm_asound"

LABEL="end"
